# crontab
0 */4 * * * root python /app/backup.py >> /var/log/cron.log 2>&1
